# systemd/dropkick.service
[Unit]
Description=Dropkick projectM visualizer
After=graphical.target sound.target pipewire.service
Wants=graphical.target

[Service]
Type=simple
EnvironmentFile=%h/.local/share/dropkick/dropkick.env
EnvironmentFile=-%h/.local/share/dropkick/dropkick.local.env
# Launch via the wrapper so LD_LIBRARY_PATH (libprojectM under ~/.local/lib) and
# the display env are set even if the graphical session didn't export them.
ExecStart=%h/.local/bin/dropkick
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
